# images/mythsuse/src/Obs/supervisor-files/db-supervisord.conf
[supervisord]
nodaemon=false
[program:mariadb-server]
command=/usr/bin/mysqld_safe --skip-syslog --timezone=%(ENV_TZ)s --datadir='/db'
